# debian/strongswan-starter.preinst
#!/bin/sh

set -e

case "$1" in
    upgrade)
        if dpkg --compare-versions "$2" 'lt' '5.2.1-5pexip1'; then
            invoke-rc.d ipsec stop || true
        fi
        ;;
    install|abort-upgrade)
        ;;
    *)
        echo "$0 called with unknown argument \`$1'" 1>&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
