# debian/strongswan-starter.preinst
#!/bin/sh

set -e

case "$1" in
    upgrade)
        invoke-rc.d ipsec stop || true
        ;;
    install|abort-upgrade)
        ;;
    *)
        echo "$0 called with unknown argument \`$1'" 1>&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
